# scripts/mobilede_configwatch.service
[Unit]
Description=Restart MobileDeBotParser on config change

[Service]
Type=oneshot
ExecStart=/var/www/mobile/scripts/mobilede_docker_restart.sh
