# scripts/mobilede_configwatch.path
[Unit]
Description=Watch configuration and exclusion files for changes MobileDe Parser

[Path]
PathModified=/var/www/mobile/parser/configuration.yaml
PathModified=/var/www/mobile/excludes_brand.csv
PathModified=/var/www/mobile/excludes_dealer.csv
PathModified=/var/www/mobile/excludes_images.csv
PathModified=/var/www/mobile/replaces.csv

[Install]
WantedBy=multi-user.target
