# scripts/mobilede_configwatch.path
[Unit]
Description=Watch configuration.yaml for changes MobileDe Parser

[Path]
PathModified=/var/www/mobile/parser/configuration.yaml

[Install]
WantedBy=multi-user.target
